# attach/conky.service
[Unit]
Description=Запуск conky.
 
[Service]
ExecStart=/home/lljk/.config/conky/conky.sh
 
[Install]
WantedBy=multi-user.target
